# etc/hw-check.cron.daily
#!/bin/sh

which hw-check > /dev/null && hw-check -wm root@localhost
